# docker/supervisord/supervisord.conf
[supervisord]
nodaemon = true

[program:sshd]
priority = 1
command = /usr/sbin/sshd -D

[program:syslog]
priority = 2
command = bash -c "sleep 30 && /etc/init.d/syslog-ng start"

[program:telegraf]
priority = 3
command = bash -c "sleep 30 && /etc/init.d/telegraf start"

[program:chronograf]
priority = 4
command = /usr/bin/chronograf -r

[program:snmpwalk]
priority = 5
command = snmpwalk -c public -v 2c 172.17.0.1

